# twitch-summon.service
[Unit]
Description=Twitch Summon System
After=network.target

[Service]
Type=notify
# Pfad zu deiner App anpassen!
WorkingDirectory=/home/[USERNAME]/TwitchSummonSystem/TwitchSummonSystem/publish
ExecStart=/usr/bin/dotnet TwitchSummonSystem.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=twitch-summon
User=[USERNAME]
Group=[USERNAME]

# Environment Variables - SECURE VERSION
# DO NOT put secrets here! Use /etc/environment or EnvironmentFile instead
Environment=ASPNETCORE_ENVIRONMENT=Production
EnvironmentFile=/home/[USERNAME]/.env

[Install]
WantedBy=multi-user.target
